# deploy/runpod/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres -D /data/postgres
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres_err.log
priority=10

[program:redis]
command=redis-server --dir /data/redis
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=20

[program:milvus]
command=milvus run standalone
environment=ETCD_USE_EMBED=true,ETCD_DATA_DIR=/data/milvus/etcd,ETCD_CONFIG_PATH=/data/milvus/configs/embedEtcd.yaml
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/milvus.log
stderr_logfile=/var/log/supervisor/milvus_err.log
priority=30

[program:agent]
command=python3 -m uvicorn src.agent.server:app --host 0.0.0.0 --port 8081
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/agent.log
stderr_logfile=/var/log/supervisor/agent_err.log
priority=100

[program:analytics]
command=python3 -m uvicorn src.analytics.server:app --host 0.0.0.0 --port 8082
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/analytics.log
stderr_logfile=/var/log/supervisor/analytics_err.log
priority=100

[program:unstructured_retriever]
command=python3 -m uvicorn src.retrievers.server:app --host 0.0.0.0 --port 8086
directory=/app
environment=EXAMPLE_PATH=src/retrievers/unstructured_data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/unstructured_retriever.log
stderr_logfile=/var/log/supervisor/unstructured_retriever_err.log
priority=100

[program:structured_retriever]
command=python3 -m uvicorn src.retrievers.server:app --host 0.0.0.0 --port 8087
directory=/app
environment=EXAMPLE_PATH=src/retrievers/structured_data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/structured_retriever.log
stderr_logfile=/var/log/supervisor/structured_retriever_err.log
priority=100

[program:api_gateway]
command=python3 -m uvicorn src.api_gateway.server:app --host 0.0.0.0 --port 9000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api_gateway.log
stderr_logfile=/var/log/supervisor/api_gateway_err.log
priority=110

[program:jupyter]
command=jupyter lab --allow-root --ip=0.0.0.0 --NotebookApp.token='' --port=8888
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/jupyter.log
stderr_logfile=/var/log/supervisor/jupyter_err.log
priority=120